# svnkit/nbsd/etc/nightly-incremental.cron
# svnkit nightly-incremental.cron : nightly subversion incremental backups
# $Id$
13 9 * * 1-6 . ~/.profile; svnkit dump_all incremental
