# svnkit/nbsd/etc/weekly-full.cron
# svnkit weekly-full.cron : weekly subversion full backups
# $Id$
13 10 * * 7 . ~/.profile; svnkit dump_all full
